# backend_laravel/python-ml-service/supervisor.conf
##############################################################################
# Supervisor — mantiene el microservicio FastAPI corriendo en el VPS
#
# Instalación:
#   sudo cp supervisor.conf /etc/supervisor/conf.d/mindra-ml.conf
#   sudo supervisorctl reread && sudo supervisorctl update
#   sudo supervisorctl start mindra-ml
##############################################################################

[program:mindra-ml]
command=/var/www/vhosts/mindra.cafined.org/venv/bin/python /var/www/vhosts/mindra.cafined.org/httpdocs/python-ml-service/main.py
directory=/var/www/vhosts/mindra.cafined.org/httpdocs/python-ml-service
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/mindra-ml.err.log
stdout_logfile=/var/log/supervisor/mindra-ml.out.log
; Run as the web-server user — never as root.
; On Plesk Ubuntu the web user is typically www-data; adjust if Plesk uses
; a domain-specific user (check with: id $(stat -c %U /var/www/vhosts/mindra.cafined.org/httpdocs))
user=www-data
environment=APP_ENV="production",PYTHON_ML_PORT="8001",ML_INFERENCE_ROOT="/var/www/vhosts/mindraback.cafined.org/httpdocs"
stopasgroup=true
killasgroup=true
