# backend_laravel/python-ml-service/supervisor.conf
##############################################################################
# Supervisor — mantiene el microservicio FastAPI corriendo en el VPS
#
# Instalación:
#   sudo cp python-ml-service/supervisor.conf /etc/supervisor/conf.d/mindra-ml.conf
#   sudo supervisorctl reread && sudo supervisorctl update
#   sudo supervisorctl start mindra-ml
#
# Ruta del proyecto en el servidor:
#   /var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/
##############################################################################

[program:mindra-ml]
command=/var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/python-ml-service/venv/bin/python main.py
directory=/var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/python-ml-service
autostart=true
autorestart=true
; Reintentar hasta 10 veces antes de darse por vencido
startretries=10
; Tiempo de espera antes de reintentar (segundos)
startsecs=10
; Si el proceso muere antes de 10s, se considera que falló al arrancar
; Si muere después, autorestart lo levanta de nuevo
stopwaitsecs=30
stderr_logfile=/var/log/supervisor/mindra-ml.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile=/var/log/supervisor/mindra-ml.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; Corre como el usuario del dominio en Plesk
user=cafined
environment=APP_ENV="production",PYTHON_ML_PORT="8001",OMP_NUM_THREADS="1",TOKENIZERS_PARALLELISM="false"
stopasgroup=true
killasgroup=true
