# backend_laravel/python-ml-service/supervisor.conf
##############################################################################
# Supervisor — mantiene el microservicio FastAPI corriendo en el VPS
#
# Instalación:
#   sudo cp python-ml-service/supervisor.conf /etc/supervisor/conf.d/mindra-ml.conf
#   sudo supervisorctl reread && sudo supervisorctl update
#   sudo supervisorctl start mindra-ml
#
# Ruta del proyecto en el servidor:
#   /var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/
##############################################################################

[program:mindra-ml]
command=/var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/python-ml-service/venv/bin/python main.py
directory=/var/www/vhosts/cafined.org/mindra.cafined.org/mindraproplus/backend_laravel/python-ml-service
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/mindra-ml.err.log
stdout_logfile=/var/log/supervisor/mindra-ml.out.log
; Corre como el usuario del dominio en Plesk (ajustar si es distinto)
user=cafined
environment=APP_ENV="production",PYTHON_ML_PORT="8001"
stopasgroup=true
killasgroup=true
